# profile.d/.gitconfig
[user]
        email = DownLord@gmail.com
        name = DownLord
[push]
        default = matching
