# profile.d/tabwidth.sh
#!/bin/bash

/usr/sbin/setterm -regtabs 2 2>>/dev/null
