# profile.d/xserver.sh
#!/bin/bash

cp /home/pi/.Xauthority ~/
export DISPLAY=:0
